# BACKEND ACTUALIZADO/radius-api.service
[Unit]
Description=RADIUS Manager API
After=network.target mariadb.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/radius-backend
EnvironmentFile=/opt/radius-backend/.env
ExecStart=/opt/radius-backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8080 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
